# layers/debian/rootfspatch/rootfs/etc/update-motd.d/10-help-text
#!/bin/sh

printf "\n"
printf " * Support:        {{ distro.support }}\n"
