# layers/debian-rootfs/rootfspatch/rootfs/etc/update-motd.d/00-header
#!/bin/sh

printf "Welcome to {{ distro.name }} (%s %s %s)\n" "$(uname -o)" "$(uname -r)" "$(uname -m)"
